# ods-signerd.service
[Unit]
Description=OpenDNSSEC signer daemon
After=syslog.target network.target ods-enforcerd

[Service]
Type=simple
PIDFile=/var/run/opendnssec/signerd.pid
EnvironmentFile=-/etc/sysconfig/ods
ExecStart=/usr/sbin/ods-signerd -d $ODS_SIGNERD_OPT

[Install]
WantedBy=multi-user.target
